# console/files/ttyS0.conf
stop on runlevel [S016]
start on runlevel [2345]

respawn
instance /dev/ttyS0
exec /sbin/agetty -L ttyS0 38400 vt102
